# supervisord.conf
[supervisord]
nodaemon=true
user=root
# Основной лог супервизора — теперь наш "Общий лог" всех событий
logfile=/app/logs/all_services.log
logfile_maxbytes=50MB
logfile_backups=10
pidfile=/var/run/supervisord.pid

# Этот блок позволяет перехватывать логи всех подпрограмм и писать их в общий файл/вывод
[eventlistener:stdout]
command=tee -a /app/logs/all_services.log
events=PROCESS_LOG
result_handler=supervisor.handlers:event_handler

# --- FASTAPI (Основное API) ---
[program:fastapi]
command=uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
# Объединяем stderr и stdout в один файл
redirect_stderr=true
stdout_logfile=/app/logs/fastapi.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

# --- ENGINE WORKER (Обработка логики) ---
[program:engine]
command=python engine_worker.py
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/engine.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

# --- CONNECTOR WORKER (Сбор данных Авито) ---
[program:connector]
command=python connector_worker.py
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/connector.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

# --- TG WORKER (Уведомления) ---
[program:tg_bot]
command=python tg_worker.py
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/tg_worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

# --- SCHEDULER (Планировщик задач) ---
[program:scheduler]
command=python scheduler.py
directory=/app
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5